# scripts/com.doing2done.digest.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.doing2done.digest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO__/scripts/run-push.sh</string>
    <string>digest</string>
  </array>
  <!-- Sunday 18:00 — the week in review, before the next one starts. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>
    <key>Hour</key><integer>18</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>__REPO__/data/push.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO__/data/push.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
